# disable-wakeup.service
[Unit]
Description=Disable XHC1 and LID0 wakeup

[Service]
Type=oneshot
ExecStart=/bin/sh -c "echo 'XHC1' > /proc/acpi/wakeup; echo 'LID0' > /proc/acpi/wakeup"

[Install]
WantedBy=multi-user.target
